# systemd/reacher@.service
[Unit]
Description=REACHER API Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
# %i is the username passed at enable time, e.g. `systemctl enable reacher@pi`
Environment="PATH=/home/%i/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=reacher
Restart=on-failure
RestartSec=5

# View logs with: journalctl -u reacher@<user> -f
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
